# infrastructure/systemd/smoldocling-docling-chunk-worker-watchdog.timer
[Unit]
Description=Run Docling chunk worker heartbeat watchdog every 5 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=5m
AccuracySec=30s
Persistent=true
Unit=smoldocling-docling-chunk-worker-watchdog.service

[Install]
WantedBy=timers.target
