# infrastructure/systemd/smoldocling-docling-chunk-worker-test.service
[Unit]
Description=PDF Docling chunk NATS worker — test broker (CPU, pdf-test)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/smoldocling/apps/pdf-test
EnvironmentFile=-/home/smoldocling/apps/pdf-test/.env
Environment=OMP_NUM_THREADS=4
ExecStart=/home/smoldocling/apps/pdf-test/venv/bin/python /home/smoldocling/apps/pdf-test/docling_chunk_worker.py
Restart=always
RestartSec=10
StandardOutput=append:/home/smoldocling/apps/pdf-test/chunk-worker-test.log
StandardError=append:/home/smoldocling/apps/pdf-test/chunk-worker-test.log

[Install]
WantedBy=default.target
